# systemd/cloud-init-local.service
[Unit]
Description=Initial cloud-init job (pre-networking)
Wants=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/cloud-init start-local
RemainAfterExit=yes

# Output needs to appear in instance console output
StandardOutput=tty

[Install]
WantedBy=multi-user.target
