# systemd/cloud-init-hotplugd.service
[Unit]
Description=cloud-init hotplug hook daemon
DefaultDependencies=yes
Wants=cloud-init-local.service
After=cloud-final.service cloud-init-hotplugd.socket
ConditionPathExists=/etc/cloud/cloud-init-hook-hotplug

[Service]
Type=simple
ExecStart=/usr/bin/cloud-init devel hotplug-hook
SyslogIdentifier=cloud-init-hotplugd

[Install]
WantedBy=cloud-init.target
