# udev/10-cloud-init-hook-hotplug.rules
SUBSYSTEM=="net|block", RUN+="/usr/lib/cloud-init/hook-hotplug"
