# sysvinit/gentoo/cloud-config
#!/sbin/openrc-run

depend() {
  after cloud-init-local
  after cloud-init
  before cloud-final
  provide cloud-config
}

start() {
  grep -q 'cloud-init=disabled' /proc/cmdline \
    && ewarn "$RC_SVCNAME is disabled via /proc/cmdline." \
    && exit 0

  cloud-init modules --mode config
  eend 0
}
