# sysvinit/gentoo/cloud-config
#!/sbin/openrc-run

depend() {
  after cloud-init-local
  after cloud-init
  before cloud-final
  provide cloud-config
}

start() {
  if grep -q 'cloud-init=disabled' /proc/cmdline; then
    ewarn "$RC_SVCNAME is disabled via /proc/cmdline."
  elif test -e /etc/cloud/cloud-init.disabled; then
    ewarn "$RC_SVCNAME is disabled via cloud-init.disabled file"
  else
    cloud-init modules --mode config
  fi
  eend 0
}
