# sysvinit/gentoo/cloud-final
#!/sbin/openrc-run

depend() {
  after cloud-config
  provide cloud-final
}

start() {
  cloud-init modules --mode final
  eend 0
}
