# sysvinit/gentoo/cloud-init-local
#!/sbin/runscript

depend() {
  after net # remove after nova-agent fix
  before cloud-init
  provide cloud-init-local
}

start() {
  cloud-init init --local
  eend 0
}
