# sysvinit/gentoo/cloud-final
#!/sbin/openrc-run

depend() {
  after cloud-config
  provide cloud-final
}

start() {
  grep -q 'cloud-init=disabled' /proc/cmdline \
    && ewarn "$RC_SVCNAME is disabled via /proc/cmdline." \
    && exit 0
  cloud-init modules --mode final
  eend 0
}
