# sysvinit/gentoo/cloud-init
#!/sbin/openrc-run
# add depends for network, dns, fs etc
depend() {
  after cloud-init-local
  after net
  before cloud-config
  provide cloud-init
}

start() {
  cloud-init init
  eend 0
}
