# sysvinit/gentoo/cloud-init
#!/sbin/openrc-run
# add depends for network, dns, fs etc
depend() {
  after cloud-init-local
  after net
  before cloud-config
  provide cloud-init
}

start() {
  grep -q 'cloud-init=disabled' /proc/cmdline \
    && ewarn "$RC_SVCNAME is disabled via /proc/cmdline." \
    && exit 0
  cloud-init init
  eend 0
}
