# sysvinit/gentoo/cloud-init-local
#!/sbin/runscript

depend() {
  after localmount
  after netmount
  before cloud-init
  provide cloud-init-local
}

start() {
  cloud-init init --local
  eend 0
}
