# sysvinit/gentoo/cloud-init-local
#!/sbin/openrc-run

depend() {
  after localmount
  before net
  before cloud-init
  provide cloud-init-local
}

start() {
  grep -q 'cloud-init=disabled' /proc/cmdline && exit 0
  cloud-init init --local
  eend 0
}
