# sysvinit/gentoo/cloud-init
#!/sbin/runscript
# add depends for network, dns, fs etc
depend() {
  after cloud-init-local
  before cloud-final
  provide cloud-init
}

start() {
  cloud-init init
  eend 0
}
