# sysvinit/openrc/cloud-init-local
#!/sbin/openrc-run

depend() {
  after localmount
  before net
  before cloud-init
  provide cloud-init-local
}

start() {
  if grep -q 'cloud-init=disabled' /proc/cmdline; then
    ewarn "$RC_SVCNAME is disabled via /proc/cmdline."
  elif test -e /etc/cloud/cloud-init.disabled; then
    ewarn "$RC_SVCNAME is disabled via cloud-init.disabled file"
  else
    cloud-init init --local
  fi

  eend 0
}
