# upstart/cat-cloud-config.conf
description "Cat the Config"

start on cloud-config
console output
task

exec cloud-init-run-module once_per_ami catconfig execute cat $CLOUD_CFG
