# upstart/cat-cloud-config.conf
description "Cat the Config"

start on cloud-config
console output
task

exec cloud-init-run-module once-per-instance catconfig execute cat $CLOUD_CFG
