# upstart/cloud-config-misc.conf
# Miscellaneous cloud-config items
#  includes enabling 'updates check'
description "Run miscellaneous cloud-config items"

start on cloud-config
console output
task

exec cloud-init-cfg config-misc
