# upstart/cloud-apt-update-upgrade.conf
# cloud-apt-update-upgrade - Update software at boot
description "Update software at boot"

start on filesystem
console output

exec cloud-init-cfg apt-update-upgrade
