# upstart/cloud-disable-ec2-metadata.conf
# cloud-disable-ec2-metadata - Disable the ec2 metadata service
description "Disable the ec2 metadata service"

start on filesystem
console output
task

exec cloud-init-cfg disable-ec2-metadata
