# upstart/cloudinit.conf
#start on (mounted MOUNTPOINT=/ and net-device-up IFACE=eth0)
# the above start-on never occurs, blocking boot. [bug 504883]
# instead, start much later in the boot process. sigh..
start on (local-filesystems and net-device-up IFACE=eth0)

task

console output

exec /usr/bin/cloud-init start
