# upstart/ec2-ebs-mounts.conf
# ec2-ebs-mounts
#
# Mount EC2 EBS mount points

description	"Populate EBS mountpoints"

start on cloud-config

console output
task

exec /usr/sbin/ec2-ebs-mount
