# upstart/ec2init.conf
start on (mounted MOUNTPOINT=/ and net-device-up IFACE=eth0)

task

console output

exec /usr/bin/ec2-init start
